# debian/vncrecordd.upstart
# vncrecordd - Record a VNC session for playback with vncdotool
#
# Logs all activity to a single file per connection


start on started networking
stop on runlevel[06]

env DEFAULTFILE=/etc/default/vncrecordd
env SERVER=localhost
env PORT=5922
env OUTPUT_DIR=/tmp

respawn

script
    if [ -f $DEFAULTFILE ]; then
        . $DEFAULTFILE
    fi

    /usr/bin/vncdotool -o $OUTPUT_DIR -s $SERVER service $PORT
end script
